# tacos-celery-beat.service.example
[Unit]
Description=TaCOS Celery Beat (定时任务调度器)
After=network.target redis.service

[Service]
Type=simple
# 指定运行用户（根据实际情况修改，通常使用 www-data 或创建专用用户）
User=www-data
Group=www-data

# 工作目录 - 修改为实际部署路径
WorkingDirectory=/var/www/tacos/tacos_backend

# 从 .env.production.local 文件读取环境变量（包含敏感信息）
# 修改为实际部署路径
EnvironmentFile=/var/www/tacos/tacos_backend/.env.production.local

# 可以在这里覆盖特定的环境变量
Environment="C_FORCE_ROOT=true"

# Celery Beat 配置
# 修改为实际的虚拟环境路径
ExecStart=/var/www/tacos/tacos_backend/.venv/bin/celery -A config beat \
    --loglevel=info \
    --schedule=/var/www/tacos/tacos_backend/celerybeat-schedule \

# 运行时目录
RuntimeDirectory=celery
RuntimeDirectoryMode=755

# 安全设置
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# 修改为实际的部署路径
ReadWritePaths=/var/www/tacos/tacos_backend /var/log/tacos

# 自动重启
Restart=always
RestartSec=10

# 优雅停止
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target
